# deploy/slyck.conf
[program:slyck]
command = /home/ubuntu/play_stuff/target/universal/stage/bin/hey    -J-Xmx512m        ; Command to start app
user = ubuntu                                                          ; User to run as
stdout_logfile = /home/ubuntu/play_app.log   ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8 
